# launchd/com.abhimehrotra.email-security-pipeline.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnvironmentVariables</key>
	<dict>
		<key>HOME</key>
		<string>/Users/speedybee</string>
		<key>PATH</key>
		<string>/opt/homebrew/bin:/opt/homebrew/sbin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
	</dict>
	<key>KeepAlive</key>
	<dict>
		<key>SuccessfulExit</key>
		<false/>
	</dict>
	<key>Label</key>
	<string>com.abhimehrotra.email-security-pipeline</string>
	<key>Nice</key>
	<integer>5</integer>
	<key>ProgramArguments</key>
	<array>
		<string>/bin/zsh</string>
		<string>/Users/speedybee/dev/email-security-pipeline/launchd/start-email-security-pipeline.sh</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>StandardErrorPath</key>
	<string>/Users/speedybee/Library/Logs/email-security-pipeline/pipeline.err</string>
	<key>StandardOutPath</key>
	<string>/Users/speedybee/Library/Logs/email-security-pipeline/pipeline.out</string>
	<key>ThrottleInterval</key>
	<integer>60</integer>
	<key>WorkingDirectory</key>
	<string>/Users/speedybee/dev/email-security-pipeline</string>
</dict>
</plist>
